# Setup_and_Install/runit/rfkill/run
#!/bin/sh

sv start dhcps || exit 1
rfkill unblock all
